# roles/wol/templates/wol.service.j2
[Unit]
Description=Configure Wake On LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ default_network_interface.stdout }} wol g

[Install]
WantedBy=basic.target
